# example.crontab
* * * * * echo "hello world! $(date)"; sleep 10; echo "goodby world! $(date)"
* * * * * LANG=C date; sleep 5; LANG=C date
