# example.crontab
# comment
RUNNER=/usr/local/bin/job-runner

* * * * * echo "hello world! $(date)"; sleep 10; echo "goodby world! $(date)"
* * * * * $RUNNER -- date
