# tests/crontab
# This is a comment
  # indented comment

*  *   *   *  * echo "hello world!"
0  0  31  12  * date

# done
